# stado/deploy/templates/wisent-host-health.timer.tmpl
[Unit]
Description=Wisent host health beacon (every 60s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=wisent-host-health.service

[Install]
WantedBy=timers.target
