# wisent_compute/deploy/templates/wisent-upgrade.service.tmpl
[Unit]
Description=Wisent Compute self-upgrade (daily)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={USER}
Group={USER}
Environment=HOME={HOME}
Environment=PATH={PATH}
# Cleanup BEFORE pip install so the upgrade has disk to work with
# (otherwise pip --upgrade hits ENOSPC on a tight disk and the
# upgrade fails repeatedly without progress).
ExecStartPre=-/bin/bash {HOME}/wisent_pre_start_cleanup.sh
ExecStart=/usr/bin/python3 -m pip install --user --break-system-packages --upgrade wisent-compute wisent wisent-tools wisent-extractors wisent-evaluators
# After upgrade, refresh the cleanup script itself from the new
# package install so improvements to the cleanup logic propagate
# automatically without ever re-rendering the .service file.
ExecStartPost=-/bin/bash -c 'src=$(/usr/bin/python3 -c "import os,wisent_compute;print(os.path.join(os.path.dirname(wisent_compute.__file__),\"deploy\",\"pre_start_cleanup.sh\"))"); [ -f "$src" ] && /bin/cp "$src" {HOME}/wisent_pre_start_cleanup.sh && /bin/chmod +x {HOME}/wisent_pre_start_cleanup.sh'
StandardOutput=append:/var/log/wisent-upgrade.log
StandardError=append:/var/log/wisent-upgrade.log
