# wisent_compute/deploy/templates/wisent-upgrade.timer.tmpl
[Unit]
Description=Wisent Compute self-upgrade daily timer

[Timer]
OnBootSec=10min
OnCalendar=daily
RandomizedDelaySec=30min
Persistent=true
Unit=wisent-upgrade.service

[Install]
WantedBy=timers.target
