# wisent_compute/deploy/templates/wisent-upgrade.service.tmpl
[Unit]
Description=Wisent Compute self-upgrade (daily)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={USER}
Group={USER}
Environment=HOME={HOME}
Environment=PATH={PATH}
ExecStart=/usr/bin/python3 -m pip install --user --break-system-packages --upgrade wisent-compute wisent wisent-tools wisent-extractors wisent-evaluators
StandardOutput=append:/var/log/wisent-upgrade.log
StandardError=append:/var/log/wisent-upgrade.log
