# deploy/host-health-beacon.timer
[Unit]
Description=Wisent host health beacon (every 60s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s

[Install]
WantedBy=timers.target
